# server/xmpp/systemd/trix-xmpp-backup.timer
[Unit]
Description=Run Trix XMPP backup daily

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
